# init/concourse-bootstrap-fly.conf
start on (started concourse-web and started concourse-worker)

exec /opt/concourse/bin/fly-bootstrap
